# bin/rc.local
kldload bktr
kldload vesa
/home/nex/bin/setupwebcam
